# infrastructure/ubuntu/systemd/tochi-backend.service
[Unit]
Description=TOCHI backend
After=network.target

[Service]
Type=simple
User=tochi
Group=tochi
WorkingDirectory=/opt/tochi/TOCHI_LEGAL_SUITE
Environment=NODE_ENV=production
Environment=PORT=4000
EnvironmentFile=/opt/tochi/TOCHI_LEGAL_SUITE/backend.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
